# scripts/DEBIAN/prerm
#!/bin/bash
systemctl stop addbook.service
systemctl disable addbook
rm /etc/systemd/system/addbook.service
